# projects/bharat_market_pulse/scheduler_example.cron
# Bharat Market Pulse scheduler (IST)
# Morning update: before market open (9:00 AM IST)
0 9 * * 1-5 cd /home/ajjames/.openclaw/workspace/projects/bharat_market_pulse && /usr/bin/env bash -lc 'source .venv/bin/activate && python report_pipeline.py >> reports/scheduler.log 2>&1'

# Closing update: shortly after market close (3:45 PM IST)
45 15 * * 1-5 cd /home/ajjames/.openclaw/workspace/projects/bharat_market_pulse && /usr/bin/env bash -lc 'source .venv/bin/activate && python report_pipeline.py >> reports/scheduler.log 2>&1'
